# Build/Podman/devops/gitea.container
[Unit]
Description=Gitea-Server
After=postgres.service
Requires=postgres.service

[Service]
Restart=always
RestartSec=10s

[Container]
NoNewPrivileges=true
UserNS=keep-id:uid=1000,gid=1000

Image=docker.gitea.com/gitea:1.24.6-rootless
ContainerName=gitea
Pod=devops.pod
Network=services.network

EnvironmentFile=%h/devops/gitea.env

Volume=%h/devops/gitea/git/:/data/git/:Z
Volume=%h/devops/gitea/gitea/:/data/gitea/:Z

[Install]
WantedBy=multi-user.target default.target
