# Build/Podman/devops/gitea.container
[Unit]
Description=Gitea-Server
After=postgres.service
Requires=postgres.service

[Service]
Restart=always
RestartSec=10s

[Container]
NoNewPrivileges=true
UserNS=keep-id:uid=1000,gid=1000

Image=docker.gitea.com/gitea:1.25.2-rootless
ContainerName=gitea
Network=services.network

EnvironmentFile=%h/devops/gitea.env

Volume=%h/devops/gitea/data/:/var/lib/gitea/:Z
Volume=%h/devops/gitea/config/:/etc/gitea/:Z

[Install]
WantedBy=multi-user.target default.target
