# Build/Podman/devops/postgres.container
[Unit]
Description=PostgreSQL DB Server

[Service]
Restart=always
RestartSec=10s

[Container]
NoNewPrivileges=true
UserNS=keep-id:uid=1000,gid=1000

Image=docker.io/library/postgres:18.1-alpine3.21
ContainerName=postgres
Pod=devops.pod
Network=services.network

Volume=%h/devops/postgres/:/var/lib/postgresql/18/docker:Z

Environment=POSTGRES_PASSWORD=1234

[Install]
WantedBy=multi-user.target default.target
